# infra/systemd/observer.service
[Unit]
Description=Observer Container Service
Documentation=https://github.com/your-repo/prj_ops
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/azureuser/app/obs_deploy
EnvironmentFile=-/home/azureuser/app/obs_deploy/.env

# 컨테이너 시작
ExecStart=/usr/bin/docker-compose up -d

# 컨테이너 중지
ExecStop=/usr/bin/docker-compose down

# 재시작 설정
Restart=on-failure
RestartSec=10s

# 재시작 루프 방지
StartLimitBurst=5
StartLimitIntervalSec=60s

# 타임아웃 설정
TimeoutStartSec=300s
TimeoutStopSec=60s

# 사용자 설정
User=azureuser
Group=azureuser

[Install]
WantedBy=multi-user.target
